# debian/tests/overrides/colord.service
[Service]
ExecStart=
ExecStart=/usr/libexec/colord --create-dummy-sensor
